# deploy/agent/macos/com.iscy.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.iscy.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/libexec/iscy-agent</string>
    <string>--backend-url</string>
    <string>https://iscy.example.org</string>
    <string>--tenant-id</string>
    <string>1</string>
    <string>--state-path</string>
    <string>/Library/Application Support/ISCY Agent/state.json</string>
    <string>--queue-dir</string>
    <string>/Library/Application Support/ISCY Agent/queue</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ISCY_AGENT_CHANNEL</key>
    <string>launchd</string>
    <key>ISCY_AGENT_QUEUE_MAX_FILES</key>
    <string>100</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/var/log/iscy-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/iscy-agent.error.log</string>
</dict>
</plist>
